# autonomy/build/andretti.service
[Unit]
Description=Andretti Service for driving the pathplanner path
StartLimitIntervalSec=0
PartOf=robot.service
After=robot.service

[Service]
Type=idle
Restart=always
RestartSec=2
WorkingDirectory=/home/robot/2020/autonomy/andretti
ExecStart=/home/robot/2020/autonomy/andretti/andretti.sh

[Install]
WantedBy=robot.service