# autonomy/build/backend.service
[Unit]
Description=Backend Service for managing the robot
StartLimitIntervalSec=0
PartOf=robot.service
After=robot.service

[Service]
Type=idle
Restart=always
RestartSec=2
WorkingDirectory=/home/pi/2020/autonomy/backend
ExecStart=./backend.sh

[Install]
WantedBy=multi-user.