# autonomy/build/backend.service
[Unit]
Description=Backend Service for managing the robot
StartLimitIntervalSec=0
PartOf=robot.service
After=robot.service

[Service]
Type=idle
Restart=always
RestartSec=2
WorkingDirectory=/home/mperry37/2020/autonomy/backend
ExecStart=/home/mperry37/2020/autonomy/backend/backend.sh

[Install]
WantedBy=robot.service