# autonomy/build/cartographer.service
[Unit]
Description=Cartographer service that runs the Cartographer module for making the map
StartLimitIntervalSec=0
PartOf=robot.service
After=robot.service

[Service]
Type=idle
Restart=always
RestartSec=2
WorkingDirectory=/home/robot/2020/autonomy/cartographer
ExecStart=/home/robot/2020/autonomy/cartographer/cartographer.sh
User=robot

[Install]
WantedBy=robot.service