# autonomy/build/localizer.service
[Unit]
Description=Localizer service for managing the robot location
StartLimitIntervalSec=0
PartOf=robot.service
After=robot.service

[Service]
Type=idle
Restart=always
RestartSec=2
WorkingDirectory=/home/robot/2020/autonomy/localizer
ExecStart=/home/robot/2020/autonomy/localizer/localizer.sh
User=robot

[Install]
WantedBy=robot.service