# autonomy/build/lunaview.service
[Unit]
Description=Lunaview service for managing the console view of the robot.
StartLimitIntervalSec=0
PartOf=robot.service
After=robot.service

[Service]
Type=idle
Restart=always
RestartSec=2
WorkingDirectory=/home/pi/2020/autonomy/lunabug
ExecStart=./lunaview.sh

[Install]
WantedBy=multi-user.