# autonomy/build/lunaview.service
[Unit]
Description=Lunaview service for managing the console view of the robot.
PartOf=robot.service
After=robot.service

[Service]
Type=idle
Restart=always
RestartSec=2
WorkingDirectory=/home/robot/2020/autonomy/lunabug
ExecStart=/home/robot/2020/autonomy/lunabug/lunaview.sh
User=robot

[Install]
WantedBy=robot.service