# autonomy/build/pathplanner.service
[Unit]
Description=Pathplanning service for managing the module of path planning
StartLimitIntervalSec=0
PartOf=robot.service
After=robot.service

[Service]
Type=idle
Restart=always
RestartSec=2
WorkingDirectory=/home/robopi/2020/autonomy/pathplanner
ExecStart=/home/robopi/2020/autonomy/pathplanner/pathplanner.sh

[Install]
WantedBy=robot.service