# autonomy/build/vision.service
[Unit]
Description=Vision service for managing the vision module for robot.
StartLimitIntervalSec=0
PartOf=robot.service
After=robot.service

[Service]
Type=idle
Restart=always
RestartSec=2
WorkingDirectory=/home/robopi/2020/autonomy/vision
ExecStart=/home/robopi/2020/autonomy/vision/vision.sh

[Install]
WantedBy=robot.service